# lib/docker-sysmd-containers.service
[Unit]
Description=Docker Containers Cleanup

[Service]
Type=oneshot
ExecStart= /bin/bash -c "docker rm --force $(docker ps -a -q)"
