# lib/docker-sysmd.service
[Unit]
Description=Docker Images Cleanup

[Service]
Type=oneshot
ExecStart=/usr/bin/docker rmi --force $(docker images -q)
