# lib/docker-sysmd.service
[Unit]
Description=Docker Images Cleanup

[Service]
Type=oneshot
ExecStart=/bin/bash -c "docker rmi -f $(docker images -q)"
