# lib/docker-sysmd.service
[Unit]
Description=Docker Images Cleanup

[Service]
Type=oneshot
ExecStart=docker rmi --force $(docker images -q)
